# config/crontab/root
@reboot /home/pi/setup/iptables-config
@reboot cd /home/pi/portal; python server.py &
